# systemd-user/bisync-personal-vault.service
[Unit]
Description=Bidirectional sync of Personal Vault between Google Drive and OneDrive
After=network-online.target rclone-onedrive.service rclone-gdrive.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync gdrive:"Personal Vault" OneDrive:"Personal-Vault" \
    --create-empty-src-dirs \
    --conflict-resolve newer \
    --conflict-loser num \
    --max-delete 25 \
    --log-file %h/.cache/bisync-personal-vault.log \
    --log-level INFO

[Install]
WantedBy=default.target
