# systemd-user/bisync-family-vault.timer
[Unit]
Description=Daily bisync of Family Vault (Drive ↔ OneDrive) — staggered 30min after Personal

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
