# systemd-user/daily-sync.timer
[Unit]
Description=Run daily sync at 10pm

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
