# systemd-user/daily-sync.service
[Unit]
Description=Daily sync of dotfiles and obsidian repos

[Service]
Type=oneshot
ExecStart=/home/laurent/scripts/daily_sync.sh
