# systemd-user/rclone-gdrive.service
[Unit]
Description=rclone mount Google Drive at ~/GoogleDrive
After=network-online.target
Wants=network-online.target
AssertPathIsDirectory=%h/GoogleDrive

[Service]
Type=notify
ExecStart=/usr/bin/rclone mount gdrive: %h/GoogleDrive \
    --vfs-cache-mode writes \
    --vfs-cache-max-size 2G \
    --dir-cache-time 1h \
    --poll-interval 15s \
    --log-level INFO \
    --log-file %h/.cache/rclone-gdrive.log
ExecStop=/usr/bin/fusermount -u %h/GoogleDrive
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
